# build/package/debian/bookbot/DEBIAN/postrm
#!/usr/bin/env sh

systemctl disable tgbot \
  && systemctl daemon-reload
